# monitor/argusreach-monitor.service
[Unit]
Description=ArgusReach Reply Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=argus
WorkingDirectory=/home/argus/.openclaw/workspace/argusreach/monitor
EnvironmentFile=/home/argus/.openclaw/workspace/argusreach/monitor/.env
ExecStart=/usr/bin/python3 /home/argus/.openclaw/workspace/argusreach/monitor/monitor.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/argus/.openclaw/workspace/argusreach/monitor/logs/monitor.log
StandardError=append:/home/argus/.openclaw/workspace/argusreach/monitor/logs/monitor.log

[Install]
WantedBy=multi-user.target
