# ops/argusreach-watcher.service
[Unit]
Description=ArgusReach Admin Auto-Restart Watcher
After=network.target argusreach-admin.service

[Service]
Type=simple
User=argus
WorkingDirectory=/home/argus/.openclaw/workspace/argusreach
ExecStart=/usr/bin/python3 /home/argus/.openclaw/workspace/argusreach/ops/admin-watcher.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
