# ops/argusreach-memory-backup.timer
[Unit]
Description=ArgusReach Memory Backup — nightly at 3am UTC

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
