# ops/argusreach-sync.timer
[Unit]
Description=ArgusReach — Instantly Sync (hourly)

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
